# docker/supervisor/loan_management_local.conf
[supervisord]
nodaemon=true

[program:roadrunner]
command=/var/www/rr serve
autostart=true
autorestart=true

